# maxbot.service
[Unit]
Description=MaxBot Auditbot — Max Messenger audit funnel bot
After=network.target

[Service]
Type=simple
User=direct
WorkingDirectory=/home/direct/MaxBot
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
